# rootfs/kconfig/IN_PLACE_UPGRADE/etc/init.d/ospkg_bootcomplete
#!/bin/sh /etc/rc.common
# {# jinja-parse #}
INSTALL_PREFIX={{INSTALL_PREFIX}}
START=99

start()
{
    $INSTALL_PREFIX/tools/ospkg boot_complete
}
